# hyperagent/schedule/archive/com.astrajax.clive-curator-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.astrajax.clive-curator-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/matthewhopkinson/Documents/AstraJax/hyperagent/scripts/run_curator_daily.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/matthewhopkinson/Documents/AstraJax/hyperagent/logs/curator-launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/matthewhopkinson/Documents/AstraJax/hyperagent/logs/curator-launchd-stderr.log</string>
</dict>
</plist>
